# cronconfig.cron
SHELL=/bin/bash
0,15,30,45 * * * * /usr/bin/python /root/gutenbergshorse/main.py >> /tmp/gutenbergshorse.log
